# code/activity_prediction/server/activity_prediction
#!/bin/sh

### BEGIN INIT INFO
# Provides:          activity_prediction
# Required-Start:
# Required-Stop:
# Should-Start:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Human activity prediction server by Mikhail Karasikov
### END INIT INFO

DAEMON_NAME=activity_prediction
DAEMON_USER=user
DIR="/home/$DAEMON_USER/activity_prediction/"
DAEMON=$DIR/activity_prediction_server.py
DAEMON_OPTS="-P 35364"


PIDFILE=/var/run/$DAEMON_NAME.pid
LOGFILE=$DIR/log

do_start () {
  echo "Starting system $DAEMON_NAME daemon"
  /sbin/start-stop-daemon --start --pidfile $PIDFILE --make-pidfile --user $DAEMON_USER --chuid $DAEMON_USER --startas $DAEMON -- start $DAEMON_OPTS 2>>"$LOGFILE" &
  echo $?
}

do_stop () {
  echo "Stopping system $DAEMON_NAME daemon"
  $DAEMON stop $DAEMON_OPTS
  echo $?
  rm $PIDFILE
}

case "$1" in
  start|stop)
    do_${1}
    ;;
  restart|reload|force-reload)
    do_stop
    do_start
    ;;
  status)
    status_of_proc "$DAEMON_NAME" "$DAEMON" && exit 0 || exit $?
    ;;
  *)
    echo "Usage: /etc/init.d/$DAEMON_NAME {start|stop|restart|status}"
    exit 1
    ;;
esac

exit 0
